# ops/com.zkyko.notification-tick.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.zkyko.notification-tick</string>
  <key>ProgramArguments</key><array>
    <string>/usr/local/bin/node</string>
    <string>/Users/zk/Trading/scripts/notification-tick.mjs</string>
    <string>tick</string>
  </array>
  <key>WorkingDirectory</key><string>/Users/zk/Trading</string>
  <key>EnvironmentVariables</key><dict>
    <key>THESIS_LOOP_URL</key><string>http://127.0.0.1:3000</string>
    <key>SYNC_SECRET</key><string>change-me</string>
  </dict>
  <!-- Every 15 minutes while the Mac is awake -->
  <key>StartInterval</key><integer>900</integer>
  <key>StandardOutPath</key><string>/tmp/zkyko-notification-tick.log</string>
  <key>StandardErrorPath</key><string>/tmp/zkyko-notification-tick-error.log</string>
</dict></plist>
